# crontab
# m h dom mon dow	command
* * * * *	/cron/scripts/example.sh >> /cron/log/example
# * * * * 7	/cron/scripts/script.sh >> /cron/log/script
#
